# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:openvpn]
command=/usr/sbin/openvpn --cd /etc/openvpn --config server.conf
